# launchd/com.example.cryptotrading-telegram.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.cryptotrading-telegram</string>

    <!-- This is a system LaunchDaemon (loaded into the `system` domain, not
         a per-user `gui/<uid>` LaunchAgent -- see the install steps in
         README.md). On at least one real Mac, a gui-domain agent's
         RunAtLoad/StartInterval never fired automatically even once across
         29+ hours and dozens of observed dark-wake cycles, while manual
         `launchctl kickstart` always worked -- GUI-domain agents tied to a
         login session were apparently not being serviced during brief
         Power Nap windows. System daemons are exactly what that mode is
         designed to run. UserName keeps it running as your normal user
         (not root) for correct file ownership on .env, telegram_state.json,
         fulldata/*.csv, etc. If gui-domain StartInterval works fine on your
         machine, a plain LaunchAgent is simpler -- try that first and only
         move to a system daemon if you hit the same silent-never-fires
         symptom. -->
    <key>UserName</key>
    <string>your-macos-username</string>

    <!-- Wrapped in caffeinate -s -i: on a laptop, macOS dark-wakes for very
         brief windows (observed as short as ~1 minute on battery) to do
         maintenance and falls right back asleep. A model-retrain check can
         take 1-2 minutes, and without this wrapper it can get killed
         mid-run before finishing or even flushing its logs -- every run
         silently vanishes and you stop getting alerts with no error
         anywhere. caffeinate holds the system awake for exactly as long as
         this one process takes, then lets it sleep again normally
         afterward. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>-i</string>
        <string>/path/to/cryptoTrading_fork/venv_telegram/bin/python3</string>
        <string>-u</string>
        <string>/path/to/cryptoTrading_fork/telegram_notify.py</string>
        <string>--once</string>
        <string>--digest-hour</string>
        <string>9</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/cryptoTrading_fork</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- launchd-native scheduling instead of the script's own --schedule
         loop: launchd tracks missed StartInterval firings (e.g. while the
         Mac was asleep) and runs the job promptly on wake, which an
         in-process scheduler inside a long-running Python process cannot
         do -- it just misses the window entirely. -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/path/to/cryptoTrading_fork/telegram_notify.out.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/cryptoTrading_fork/telegram_notify.err.log</string>
</dict>
</plist>
